# docs/systemd/kyrodb-engine.service
[Unit]
Description=KyroDB Engine
After=network.target

[Service]
ExecStart=/usr/local/bin/kyrodb-engine serve 0.0.0.0 3030 --data-dir /var/lib/kyrodb --wal-segment-bytes 134217728 --wal-max-segments 8 --compact-interval-secs 300 --compact-when-wal-bytes 536870912
WorkingDirectory=/var/lib/kyrodb
Restart=always
RestartSec=3
LimitNOFILE=1048576
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
